# etc/supervisor.d/consumer.ini
[program:weblive]
directory=/devtools/jenkins/workspace/web_live/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/web_live/target/qmzb-web-live.jar --spring.profiles.active=test --server.port=3081
stdout_logfile=/opt/logs/web-live.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true



[program:weblogin]
directory=/devtools/jenkins/workspace/web_login/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/web_login/target/qmzb-web-login.jar --spring.profiles.active=test --server.port=3082
stdout_logfile=/opt/logs/web-login.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
[root@app2 supervisord.d]# cat consumer.ini.bak 
[program:admin]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-web-admin/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-web-admin/target/qmzb-web-admin.jar --spring.profiles.active=test --server.port=7081
stdout_logfile=/opt/logs/qmzb-web-admin.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true

[program:app]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-web-app/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-web-app/target/qmzb-web-app.jar --spring.profiles.active=test --server.port=7082
stdout_logfile=/opt/logs/qmzb-web-app.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true


[program:family]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-web-family/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-web-family/target/qmzb-web-family.jar --spring.profiles.active=test --server.port=7083
stdout_logfile=/opt/logs/qmzb-web-family.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true


[program:job]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-web-job/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-web-job/target/qmzb-web-job.jar --spring.profiles.active=test --server.port=7084
stdout_logfile=/opt/logs/qmzb-web-job.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true


[program:wechat]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-web-wechat/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-web-wechat/target/qmzb-web-wechat.jar --spring.profiles.active=test --server.port=70845
stdout_logfile=/opt/logs/qmzb-web-wechat.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true


[program:queue]
directory=/devtools/jenkins/workspace/qmzb_test/qmzb-app-queue/target
command=/usr/java/jdk1.8.0_171/bin/java -Xmx400m -Xms400m -verbosegclog:log/jvm.log -XX:MaxTenuringThreshold=5 -XX:-UseAdaptiveSizePolicy -XX:SurvivorRatio=3 -XX:NewRatio=1 -XX:+PrintGCDateStamps -XX:-UseConcMarkSweepGC -XX:-UseParallelGC -XX:+PrintGCDetails -XX:+UseParallelGC  -jar /devtools/jenkins/workspace/qmzb_test/qmzb-app-queue/target/qmzb-app-queue.jar --spring.profiles.active=test --server.port=7086
stdout_logfile=/opt/logs/qmzb-web-queue.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
